# mvs.service.j2
[Unit]
Description=MVS 3.8j TK5 (Hercules Mainframe Emulator)
After=network.target

[Service]
Type=simple
User={{ mvs_user }}
WorkingDirectory={{ mvs_install_dir }}
ExecStart={{ mvs_install_dir }}/mvs
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mvs-tk5
Environment="REP101A=specific"
Environment="MAINSIZE=32"

[Install]
WantedBy=multi-user.target
